# dist/maint/mastodon-cleanup.timer
[Unit]
Description=Mastodon Cleanup timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target