# dist/mastodon-pg-backup.timer
[Unit]
Description=daily backups for pg
Requires=mastodon-pg-backup.service

[Timer]
Unit=mastodon-pg-backup.service
OnCalendar=daily

[Install]
WantedBy=timers.target
